# deploy/systemd/gmtk-jam-update.timer
[Unit]
Description=Update hosted GMTK jam branches every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=1s
Persistent=true
Unit=gmtk-jam-update.service

[Install]
WantedBy=timers.target
